# config/supervisor.conf

[unix_http_server]
file=/var/run/supervisor-dev.sock            ; (the path to the socket file)
chmod=0766                                   ; socket file mode (default 0700)

[inet_http_server]
port=localhost:9111
username=supervisor
password=rosivrepus

[supervisord]
logfile=/var/log/supervisor/dev-supervisord-log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB                           ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10                              ; (num of main logfile rotation backups;default 10)
loglevel=info                                   ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/dev-supervisord.pid                ; (supervisord pidfile;default supervisord.pid)
nodaemon=false                                  ; (start in foreground if true;default false)
minfds=1024                                     ; (min. avail startup file descriptors;default 1024)
minprocs=200                                    ; (min. avail process descriptors;default 200)
childlogdir=/var/log/supervisor                 ; ('AUTO' child log dir, default $TEMP)
directory=/var/www/frisk_migration              ; (Working directory for supervisord spawned processes)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor-dev.sock

[group:all]
programs=worker.finish_locker,worker.migrate_item,worker.seed_locker

[program:worker.finish_locker]
command=bin/console frisk:worker:finish_locker
process_name=%(program_name)s.%(process_num)02d
user=www-data
numprocs=1
autostart=true
autorestart=true
killasgroup=true
stdout_logfile=/var/log/frisk-migration/workers/finish_locker.out
stderr_logfile=/var/log/frisk-migration/workers/finish_locker.err

[program:worker.migrate_item]
command=bin/console frisk:worker:migrate_item
process_name=%(program_name)s.%(process_num)02d
user=www-data
numprocs=1
autostart=true
autorestart=true
killasgroup=true
stdout_logfile=/var/log/frisk-migration/workers/migrate_item.out
stderr_logfile=/var/log/frisk-migration/workers/migrate_item.err

[program:worker.seed_locker]
command=bin/console frisk:worker:seed_locker
process_name=%(program_name)s.%(process_num)02d
user=www-data
numprocs=1
autostart=true
autorestart=true
killasgroup=true
stdout_logfile=/var/log/frisk-migration/workers/seed_locker.out
stderr_logfile=/var/log/frisk-migration/workers/seed_locker.err
